# supervisord-carbon.conf
[supervisord]
nodaemon=true

[unix_http_server]
file = /var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:graphite-api]
command = gunicorn -w4 graphite_api.app:app -b 0.0.0.0:8888
autostart = true
autorestart = true

[program:carbon-cache]
command = /opt/graphite/bin/carbon-cache.py --nodaemon --pidfile=/var/run/carbon-cache.pid start
autostart = true
autorestart = true
